# systemd/plasmalogin.service.d/10-powos-wantedby.conf
# PowOS: force plasmalogin.service into graphical.target's Wants set.
#
# Upstream plasmalogin.service ships with [Install] Alias=display-manager.service
# and NO WantedBy=graphical.target. That means `systemctl enable` only adds the
# alias symlink — not a graphical.target.wants entry — so graphical.target
# reaching ACTIVE doesn't actually pull plasmalogin in.
#
# Fedora Kinoite ships fedora-kinoite-plasmalogin-workaround.service to fix
# this, but it's not guaranteed to exist in every Bazzite/nvidia-open variant.
# This drop-in is the belt-and-braces: adds a proper WantedBy so the
# subsequent `systemctl enable plasmalogin.service` in the Containerfile does
# create /etc/systemd/system/graphical.target.wants/plasmalogin.service via
# the normal install path, guaranteed on every variant.

[Install]
WantedBy=graphical.target
