# systemd/powos-firstboot-disk.service
# powos-firstboot-disk.service - complete the disk on the FIRST boot from a
# freshly flashed USB: create the POWOS-DATA persistence partition in the free
# space after the raw image and add the Install/Recovery boot-menu entries.
#
# WHY: baking POWOS-DATA + boot entries into the raw at BUILD time needs loop
# devices, whose partition nodes are unreliable in CI/Docker. Instead the raw
# ships as a plain bootable OS and self-completes here, on real hardware where
# partition scanning works. The wrapper (powos-firstboot-disk) invokes
# install-to-usb.sh --self-complete, then drops a marker so this runs ONCE.
#
# Self-disabling: ConditionPathExists=! on the marker means once persistence is
# set up the unit is skipped on every subsequent boot. Every step is
# best-effort — the OS already booted, so this only ADDS persistence and must
# never fail the boot.

[Unit]
Description=PowOS first-boot disk self-completion (create POWOS-DATA + boot entries)
Documentation=man:powos-firstboot-disk
ConditionPathExists=!/var/lib/powos/firstboot-disk-done
# Only fire on the USB live-boot path. On an INSTALLED system there is no
# raw-image tail to expand, no need for POWOS-DATA persistence (the disk IS
# persistent), and running install-to-usb.sh --self-complete against the boot
# NVMe would repartition the wrong device. rd.powos.ramboot=1 is set only by
# the live USB image's dracut module; installed deployments don't carry it.
ConditionKernelCommandLine=rd.powos.ramboot
# Needs all local filesystems (incl. /boot and /boot/efi) mounted so we can
# resolve the boot disk and write BLS entries to the live loader dir.
After=local-fs.target
# Run before the user gets a session, but it is not required by those targets.
Before=graphical.target multi-user.target

[Service]
Type=oneshot
RemainAfterExit=yes
# Mirror the wrapper's step-by-step progress to the SCREEN as well as the
# journal: this is the "preparing system, may take a while" first-boot path and
# a silent console reads as a hang. journal+console echoes each [N/4] line.
StandardOutput=journal+console
StandardError=journal+console
ExecStart=/usr/bin/powos-firstboot-disk
# The wrapper always exits 0; belt-and-suspenders so a nonzero can never fail
# the boot transaction.
SuccessExitStatus=0 1

[Install]
WantedBy=multi-user.target
