# systemd/powos-usb-monitor.service
[Unit]
Description=PowOS USB Hotplug Monitor
Documentation=man:homefs(8)
After=systemd-udevd.service

[Service]
Type=simple
ExecStart=/usr/local/bin/homefs-usb-monitor
Restart=always
RestartSec=5s

# Security settings
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=homefs-usb

[Install]
WantedBy=multi-user.target
