# systemd/powos-ramboot-init.service
[Unit]
Description=PowOS RAM Boot Initialization
After=local-fs.target
Before=display-manager.service
ConditionPathExists=/run/powos/ramboot-state

[Service]
Type=oneshot
# NOTE: /usr/lib/powos/overlay-init.sh only exists inside the initramfs
# (module-setup.sh inst_simple). On the real rootfs the script lives in the
# dracut module directory installed by the Containerfile.
ExecStart=/usr/lib/dracut/modules.d/90powos-ramboot/powos-overlay-init.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
