# systemd/powos-installer.service
[Unit]
Description=PowOS Interactive Installer (boot-menu "Install PowOS" entry)
Documentation=man:powos(8)
# Only runs when the USB boot menu's "Install PowOS" entry was chosen, which
# appends `powos.install=1` to the kernel command line. A normal "PowOS Live"
# boot has no such arg, so this service stays inactive.
ConditionKernelCommandLine=powos.install
# On a KDE image graphical.target starts SDDM, which seizes tty1 and hides the
# installer. Conflict with + order after the display manager so it is stopped
# before we take the console. (Belt and braces: the BLS install entry also
# forces systemd.unit=multi-user.target so the DM never starts.)
After=multi-user.target display-manager.service
Conflicts=getty@tty1.service display-manager.service

[Service]
Type=idle
# Take over tty1 for the guided installer (TUI, or the GUI when a desktop is
# up). The wizard collects disk/sizes/GPU/hostname/user/ssh/ramboot/AI/restore,
# then drives `powos install-system` + firstboot config. (The raw
# `powos install-system` is still available by hand.)
ExecStart=/usr/bin/powos-install-wizard
StandardInput=tty
StandardOutput=tty
StandardError=journal
TTYPath=/dev/tty1
TTYReset=yes
TTYVHangup=yes
# If the user quits the installer, drop back to a login rather than looping.
Restart=no

[Install]
WantedBy=multi-user.target
