# systemd/powos-cachefs-sync.service
[Unit]
Description=PowOS CacheFS Sync Daemon
Documentation=man:powos(8)
# After= on a missing unit is harmless; Wants= is NOT — it created a dangling
# "not-found" powos-homefs ref on every install, because that unit is never
# shipped (the homefs binary doesn't exist yet; CacheFS is experimental).
# Re-add Wants= only when homefs actually ships.
After=powos-homefs.service local-fs.target
# Same guard layer-sync uses: this daemon syncs the USB overlay model's home.
# Without it, fixed (bootc) installs ran the daemon forever, pointlessly
# polling a /mnt/powos-usb path that never appears.
ConditionPathExists=/run/powos/layer-paths

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/lib/powos/cachefs/cachefs-sync.py
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=10

# Security hardening
NoNewPrivileges=false
ProtectSystem=false
PrivateTmp=false

# Environment
Environment=PYTHONUNBUFFERED=1

# Graceful shutdown - sync dirty files before stopping
ExecStop=/usr/bin/python3 /usr/lib/powos/cachefs/cachefs-sync.py --sync-now
TimeoutStopSec=120

[Install]
WantedBy=multi-user.target
