# systemd/powos-cachefs-sync.service
[Unit]
Description=PowOS CacheFS Sync Daemon
Documentation=man:powos(8)
# homefs was removed entirely (dead unshipped generation), so no homefs refs.
After=local-fs.target
# Same guard layer-sync uses: this daemon syncs the USB overlay model's home.
# Without it, fixed (bootc) installs ran the daemon forever, pointlessly
# polling a USB path that never appears.
ConditionPathExists=/run/powos/layer-paths
# NOTE: CacheFS write-back is NOT implemented — this daemon currently finds no
# dirty files to sync (the FUSE process and this daemon disagree on the cache
# format). CacheFS must stay disabled until that is rebuilt; see CLAUDE.md.

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/lib/powos/cachefs/cachefs-sync.py
# No ExecReload: the daemon handles only SIGTERM/SIGINT — SIGHUP would kill it.
Restart=on-failure
RestartSec=10

# Security hardening
NoNewPrivileges=false
ProtectSystem=false
PrivateTmp=false

# Environment
Environment=PYTHONUNBUFFERED=1

# No ExecStop --sync-now: the daemon's argparse has no such flag (it exited 2
# on every stop). Shutdown is SIGTERM, handled by the daemon.
TimeoutStopSec=120

[Install]
WantedBy=multi-user.target
