# systemd/powos-homefs-sync.service
[Unit]
Description=PowOS HomeFS Sync Daemon
Documentation=man:homefs(8)
After=powos-homefs.service
Requires=powos-homefs.service

[Service]
Type=simple
ExecStart=/usr/local/bin/homefs-sync-daemon
Restart=on-failure
RestartSec=5s

# Security settings
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=homefs-sync

[Install]
WantedBy=multi-user.target
